# scripts/files-7628/etc/uci-defaults/99-set-country-cn
#!/bin/sh

uci set wireless.radio0.country='CN'
uci commit wireless

exit 0
